# gentoo/util-vserver
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

start() {
	: ${UTIL_VSERVER_VARS:=/usr/lib/util-vserver/util-vserver-vars}
	if [ ! -e ${UTIL_VSERVER_VARS} ]; then
		eerror "Cannot find util-vserver installation"
		eerror "(the file '$UTIL_VSERVER_VARS' would be expected)"
		return 1
	fi
	
	. ${UTIL_VSERVER_VARS}

	ebegin "Creating directories for $PACKAGE_NAME"
    $__PKGLIBDIR/bash-wrapper 'create_dirs'
	eend $?

	ebegin "Setting vshelper path to $_VSHELPER"
    $__PKGLIBDIR/bash-wrapper 'set_helper'
	eend $?

	ebegin "Loading default device map"
	$__PKGLIBDIR/bash-wrapper 'handleDeviceMap --set 0 "$__CONFDIR/.defaults/apps/vdevmap"'
	eend $?

	$__PKGLIBDIR/bash-wrapper 'hasCgroup'
	if [ $? -eq 0 ]; then
		ebegin "Mounting cgroup-hierarchy"
		$__PKGLIBDIR/bash-wrapper 'mount_cgroup'
		eend $?
	fi
}

stop() {
	: ${UTIL_VSERVER_VARS:=/usr/lib/util-vserver/util-vserver-vars}
	if [ ! -e ${UTIL_VSERVER_VARS} ]; then
		eerror "Cannot find util-vserver installation"
		eerror "(the file '$UTIL_VSERVER_VARS' would be expected)"
		return 1
	fi
	
	. ${UTIL_VSERVER_VARS}

	ebegin "Stopping all running guests"
	$_START_VSERVERS -j 1 --all --stop
	eend $?

	ebegin "Killing all running contexts"
	$__PKGLIBDIR/bash-wrapper 'kill_contexts'
	eend $?

	$__PKGLIBDIR/bash-wrapper 'hasCgroup'
	if [ $? -eq 0 ]; then
		ebegin "Unmounting cgroup-hierarchy"
		$__PKGLIBDIR/bash-wrapper 'umount_cgroup'
		eend $?
	fi
}

# vim:ts=4:filetype=gentoo-init-d
