# util-vserver/sysv/v_gated.subst
#!/bin/sh
# chkconfig: - 32 75
# description: Wrapper to start gated bound to a single IP

# normally, you will need /etc/vservices/gated.conf
# with the following line
# IP=0.0.0.0
# so gated is started without restriction even if you
# restart it from an ssh session, itself "chbounded" to a single
# IP.
USR_LIB_VSERVER=/usr/lib/util-vserver
exec $USR_LIB_VSERVER/vsysvwrapper gated $*

