# util-vserver/sysv/v_sshd.subst
#!/bin/sh
# chkconfig: - 55 25
# description: Wrapper to start sshd bound to a single IP
USR_LIB_VSERVER=/usr/lib/util-vserver
exec $USR_LIB_VSERVER/vsysvwrapper sshd $*






