# util-vserver/sysv/v_xinetd.subst
#!/bin/sh
# chkconfig: - 56 50
# description: Wrapper to start xinetd bound to a single IP
USR_LIB_VSERVER=/usr/lib/util-vserver
exec $USR_LIB_VSERVER/vsysvwrapper xinetd $*

