# util-vserver/sysv/vkillall
#!/bin/sh
# chkconfig: - 27 73
# description: Stops unhandled vservers



