# systemd/vserver.target
# This file is part of util-vserver

[Unit]
Description = vserver

[Install]
WantedBy=multi-user.target
