# systemd/vserver-mark@.target
[Unit]
Description = Target for vservers tagged with '%I'

[Install]
WantedBy = vserver.target
