# debian/util-vserver-legacy.postrm
#!/bin/sh

if test "$1" = purge; then
    update-rc.d rebootmgr remove
    update-rc.d vservers-legacy remove
    for i in httpd named portmap sendmail smb sshd xinetd gated; do
       update-rc.d v_$i remove
    done
fi
