# distrib/gentoo/net.vserver
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
        provide net
}

start() {
        return 0
}

stop() {
        return 0
}
